# services/mnt-nfs.automount
[Unit]
Description=Automount /mnt/nfs

[Automount]
Where=/mnt/nfs

[Install]
WantedBy=multi-user.target
